# src/RaspberryPI/BrokerConnector.conf
[program:brokerConnector]
command= python3 main.py
directory= /var/brokerConnector
user=pi
autostart=true
autorestart=true
